# homecaptain/homecaptain/celerybeatd.conf
[program:celerybeat]
directory=/home/app/home-captain-dev/homecaptain
command=/home/app/home-captain-dev/env/bin/celery beat -A homecaptain --schedule=celerybeat-schedule --loglevel=INFO --pidfile=/run/homecaptain/celerybeat.pid

user=app
numprocs=1
stdout_logfile=/var/log/homecaptain/beat.log
stderr_logfile=/var/log/homecaptain/beat.log
autostart=true
autorestart=true
startsecs=10

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

priority=999